# crontab
# Set PATH for cron environment
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# Run metering processor every 5 minutes (with lock to prevent overlapping)
*/5 * * * * cd /app && set -a && source /app/cron.env && set +a && flock -n /tmp/metering_processor.lock /usr/local/bin/python3 src/metering_processor.py >> /var/log/cron.log 2>&1

# Empty line required at end of crontab
